# spec/rendered_templates/storm-nimbus.conf
description "Storm nimbus upstart"
author  "Kai Sasaki <lewuathe@me.com>"

start on runlevel [2345]
stop on runlevel [016]

chdir /usr/share/storm/0.9.3
exec bin/storm nimbus 2>&1
respawn
